# common_AOK/hostname_handling/aok-hostname-service
#!/sbin/openrc-run
#
#  Part of https://github.com/jaclu/AOK-Filesystem-Tools
#
#  Copyright (c) 2023: Jacob.Lundqvist@gmail.com
#
#  License: MIT
#
#  Appends -aok to hostname if the kernel is iSH-AOK
#  This is only installed if AOK_HOSTNAME_SUFFIX was "Y"
#

hn_alt_bin="/usr/local/bin/hostname"
if [ -x "$hn_alt_bin" ]; then
    hn_bin="$hn_alt_bin"
else
    hn_bin=/bin/hostname
fi

hn_current="$($hn_bin -s)"
#  file alternate hostname reads to find hostname
hn_source_fname="/etc/opt/hn_source_fname"

description="Ensures hostname has -aok suffix for iSH-AOK kernels."

log_file=/var/log/services.log
log_it() {
    # if no supported loggers are found, fall back to a non-dependency solution
    f_fake_syslog=/usr/local/bin/fake_syslog
    if [ -x "$f_fake_syslog" ]; then
        FAKE_SYSLOG_FILE="$log_file" "$f_fake_syslog" "$RC_SVCNAME" "$*"
    else
        echo "$(date +"%Y-%m-%d %H:%M:%S") ${RC_SVCNAME}: $*" >>/var/log/messages
    fi
}

depend() {
    if [ -f /etc/init.d/runbg ]; then
        after runbg
    fi
}

start_post() {
    log_it started
}

stop_post() {
    log_it stopped
}

start() {
    #
    #  This init script should not be present
    #  on a regular iSH, but if an FS was imported it could in
    #  theory happen...
    #
    hn_suffix="aok"

    #  If hostname -s allready ends with $hn_suffix, nothing needs to be done
    echo "$hn_current" | grep -q "\-$hn_suffix" && return

    #  Only continue if this is iSH-AOK
    grep -qi aok /proc/ish/version || return

    echo "$hn_current-$hn_suffix" >/etc/hostname

    #
    #  iSH-AOK can set hostname, so once this is done things hardcoded
    #  to use the default /bin/hostname will use the intended hostname
    #
    /bin/hostname -F /etc/hostname

    #  Since hostname might have changed, ensure it is in /etc/hosts
    /usr/local/sbin/ensure_hostname_in_host_file.sh

    /usr/local/bin/wall -n "hostname is now: $($hn_bin)"

    eend 0 # we have gotten this far, so it should be safe to signal success
}
