# common_AOK/aok_hostname/aok-hostname-service
#!/sbin/openrc-run

#
#  Needed on iSH-AOK to use my hostname convention with -aok suffix
#
description="Sets the -aok hostname of the machine."

depend() {
    if [ -f /etc/init.d/runbg ]; then
        after runbg
    fi
}

start_post() {
    echo "[$(date)] $RC_SVCNAME started" >>/var/log/services.log
}

stop_post() {
    echo "[$(date)] $RC_SVCNAME stopped" >>/var/log/services.log
}

start() {
    #
    #  This init script should not be present
    #  on a regular iSH, but if an FS was imported it could in
    #  theory happen...
    #
    if grep -qi aok /proc/ish/version; then
        #
        #  Since it is not needed for iSH-AOK, remove AOK FS custom
        #  hostname, it is only run during bootup to get the current
        #  hostname, thus it is run via /opt/AOK/...
        #  Once retrieved the normal hostname is able to use this.
        #
        rm -f /usr/local/bin/hostname
    fi

    h_suffix="aok"

    # #  If this file is present, assume iOS 17 hostname workarround should be used
    # if [ -f /etc/opt/hostname_sync_fname ]; then
    #     #
    #     #
    #     #  Since current hostname is provided via a Shortcut run as app
    #     #  starts, always use this without caching to ensure that if iOS
    #     #  devis is given a new name, it is picked up
    #     #
    #     hostname_by_shortcut="$(/opt/AOK/common_AOK/usr_local_bin/hostname --update)"
    #     echo "$hostname_by_shortcut" | grep -q "\-$h_suffix" || {
    #         echo "${hostname_by_shortcut}-$h_suffix" >/etc/hostname
    #     }
    # else

    hostname | grep -q "\-$h_suffix" || {
        echo "$(hostname)-$h_suffix" >/etc/hostname
    }

    ebegin "Setting hostname"
    hostname -F /etc/hostname
    /usr/local/bin/wall -n "hostname is now: $(hostname)"

    #
    #  This one is run from /etc/inittab, so most likely before this
    #  has started, running it again here ensures that the now active
    #  hostname is in /etc/hosts
    #
    /usr/local/sbin/ensure_hostname_in_host_file.sh
    eend $?
}
