# common_AOK/aok_hostname/aok-hostname-service
#!/sbin/openrc-run

#
#  Needed on iSH-AOK to use my hostname convention with -aok suffix
#
description="Sets the -aok hostname of the machine."

depend() {
	keyword -prefix -lxc -docker
}

start() {
	echo "NEW_HOSTNAME" >/etc/hostname
	ebegin "Setting hostname"
	hostname -F /etc/hostname
	wall "hostname is now: $(hostname)"
	#
	#  This one is run from /etc/inittab, so most likely before this
	#  has started, running it again here ensures that the now active
	#  hostname is in /etc/hosts
	#
	/usr/local/sbin/ensure_hostname_in_host_file.sh
	eend $?
}
