# common_AOK/aok_hostname/aok-hostname-service
#!/sbin/openrc-run

#
#  Needed on iSH-AOK to use my hostname convention with -aok suffix
#
description="Sets the -aok hostname of the machine."

depend() {
	keyword -prefix -lxc -docker
}

start() {
	h_suffix="aok"
	#  Ensure suffix is not added multiple times if this is restarted
	hostname | grep -qv "\-$h_suffix" && echo "$(hostname)-$h_suffix" >/etc/hostname
	ebegin "Setting hostname"
	hostname -F /etc/hostname
	/usr/local/bin/wall "hostname is now: $(hostname)"
	#
	#  This one is run from /etc/inittab, so most likely before this
	#  has started, running it again here ensures that the now active
	#  hostname is in /etc/hosts
	#
	/usr/local/sbin/ensure_hostname_in_host_file.sh
	eend $?
}
