# Debian/etc/update-motd.d/51-env-versions
#!/bin/sh

deb_version="/etc/debian_version"
aok_release="/etc/aok-release"


if [ -e "$deb_version" ]; then
    echo "Debian release: $(cat $deb_version)"
fi

if [ -e "$aok_release" ]; then
    echo "AOK release:    $(cat /etc/aok-release)"
fi

if [ -d /proc/ish ]; then
    # Dont do this whilst chrooted in linux etc
    echo "iSH release:    $(cat /proc/ish/version)"
else
    echo "iSH release:    Can't be detected whilst chrooted"
fi
