# Alpine/etc/inittab
# /etc/inittab

# Ensure all devices  are ok
::sysinit:/usr/local/sbin/fix_dev

#
#  Not system critical but we might as well do it here
#  instead of creating a full service for the task
#
::sysinit:/usr/local/sbin/ensure_hostname_in_host_file.sh

#
#  Update motd (in case iSH has been updated)
#  On iSH we get console before this is executed despite the fact that
#  sysinit actions are supposed to happen before conssole is activated
#  This means first boot after new iSH will show outdated version
#  on bootup
#
::sysinit:/usr/local/sbin/update_motd

#
#  iSH typically does no cleanup at shutdown, it just suddenly terminates.
#  When starting up openrc in default it often finds state of services as
#  already running thus leaving them as is, ie wont try to start them at all
#  since according to the states openrc finds is that default is already
#  the current runlevel.
#  Workaround is to toggle between runlevels.
#  This leads to services first being stoped and then started, ensuring
#  they start up regardless of logged runstate at startup.
#
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

#
#  I highly doubt these states can be triggered in iSH
#
# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/sbin/openrc shutdown
