# common_AOK/etc/init.d/bat-monitord
#!/sbin/openrc-run
#
#  Part of https://github.com/jaclu/AOK-Filesystem-Tools
#
#  Copyright (c) 2023: Jacob.Lundqvist@gmail.com
#
#  License: MIT
#
#  Logs Battery charge, each time it changes
#

description="Logs Battery charge, each time it changes"

command="/usr/local/sbin/bat-monitord"
command_args="> /dev/null"
command_background="YES"

pidfile="/run/bat-monitord.pid"

log_file=/var/log/services.log


start_post() {
    LOGGER_FILE="$log_file" /usr/local/bin/logger "$RC_SVCNAME" started
}

stop_post() {
    LOGGER_FILE="$log_file" /usr/local/bin/logger "$RC_SVCNAME" stopped
}
